# tools/etc/bench_agent.ini
[program:isucon]
command=/home/isucon/isucon/tools/etc/agent.sh
process_name=isucon bench agent
user=isucon
stdout_logfile=/tmp/isucon.agent.log
stderr_logfile=/tmp/isucon.agent.log
