# tools/etc/bench_agent.ini
[program:agent]
command=/Users/tagomoris/Documents/isucon/tools/etc/agent.sh
process_name=isucon bench agent
user=tagomoris
stdout_logfile=/tmp/isucon.agent.log
stderr_logfile=/tmp/isucon.agent.log
