# tools/etc/bench_master.ini
[program:master]
command=/home/isucon/isucon/tools/etc/master.sh
process_name=isucon bench master
user=isucon
stdout_logfile=/tmp/isucon.master.log
stderr_logfile=/tmp/isucon.master.log
