# tools/etc/bench_master.ini
[program:master]
command=/Users/tagomoris/Documents/isucon/tools/etc/master.sh
process_name=isucon bench master
user=tagomoris
stdout_logfile=/tmp/isucon.master.log
stderr_logfile=/tmp/isucon.master.log
