# webapp/perl/etc/isucon.ini
[program:isucon]
directory = /home/isucon/Documents/isucon/webapp/perl
command=perl -Mlib=/home/isucon/Documents/isucon/webapp/perl/extlib/lib/perl5 /home/isucon/Documents/isucon/webapp/perl/extlib/bin/plackup -s Starman -E production --preload-app --disable-keepalive --workers 10 /home/isucon/Documents/isucon/webapp/perl/app.psgi
user=isucon
stdout_logfile=/tmp/isucon.perl.log
stderr_logfile=/tmp/isucon.perl.log
