# webapp/perl/etc/isucon.ini
[program:isucon]
command=/home/isucon/isucon/webapp/perl/etc/app.sh
process_name=isucon server (perl)
user=isucon
stdout_logfile=/tmp/isucon.perl.log
stderr_logfile=/tmp/isucon.perl.log
