# webapp/nodejs/etc/isucon.ini
[program:isucon]
command=/Users/tagomoris/Documents/isucon/webapp/nodejs/etc/app.sh
process_name=isucon server (nodejs)
user=tagomoris
stdout_logfile=/tmp/isucon.nodejs.log
stderr_logfile=/tmp/isucon.nodejs.log
